# sync_wrfout_watcher.path
# ~/CODES/RUN_plots/sync_wrfout_watcher.path
[Unit]
Description=Watch for modifications to wrfout_watcher.service

[Path]
PathModified=%h/CODES/RUN_plots/wrfout_watcher.service
Unit=sync_wrfout_watcher.service

[Install]
WantedBy=default.target
